# modules/main/src/main/env/service/log/run
#!/bin/sh

LOG=${log.dir}
DAEMON=${application.daemon.enable}

if [ "$DAEMON" = "true" ]; then
  run_as="chpst -u ${application.user}:${application.group}"
fi

exec $run_as svlogd -tt \$LOG